# build/Agent/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
loglevel=debug
[program:netopeer2-server]
command=/opt/dev/Netopeer2/build/netopeer2-server -d -v %(ENV_DEBUG_LVL_NETOPEER)s
autorestart=true
redirect_stderr=true
priority=3
[program:cfgipsec2]
command=/cfgipsec2/build/i2nsf_cfgipsec2 -v %(ENV_DEBUG_LVL_CCIPS)s
autorestart=true
redirect_stderr=true
priority=4
[program:show-ip]
command=/bin/bash -c "ip a"
autorestart=false
startretries=1
priority=1
stdout_logfile=/var/log/show-ip.log
stdout_logfile_maxbytes=0
