# build/Netopeer2/supervisord.conf

[supervisord]
nodaemon=true

[program:netopeer2-server]
command=/opt/dev/Netopeer2/build/netopeer2-server -d -v ${ENV_DEBUG_LVL_NETOPEER}
environment=ENV_DEBUG_LVL_NETOPEER=2
stdout_logfile=/var/log/netopeer2-server.log
stderr_logfile=/var/log/netopeer2-server_err.log
autostart=true
autorestart=true
